# cron/etc/cron.d/memcached-restart
MAILTO=root

58 03 * * * root systemctl restart memcached.service && systemctl restart httpd.service
