# cron/etc/cron.d/memcached-restart
MAILTO=root

00 04 * * * root sudo systemctl restart memcached.service && sudo systemctl restart httpd.service
