# cron/etc/cron.d/memcached-restart
MAILTO=root

00 04 * * * root systemctl restart memcached.service && systemctl restart httpd.service
