# cron/etc/cron.d/postgres-backup
30 05 * * * postgres cd / && /usr/bin/pg_dumpall -c > /var/lib/postgresql/pg_dump/all.dump && gzip -f /var/lib/postgresql/pg_dump/all.dump 1>/dev/null
