# cron/etc/cron.d/shutdown
MAILTO=root

50 02 * * 6 root test -f /etc/cron.config/always_power_on || /sbin/shutdown -h now
