# cron/etc/cron.d/shutdown
MAILTO=root

00 22 * * 5 root /sbin/shutdown -h now
