# cron/etc/cron.d/stocks-backup
MAILTO=root

45 02 * * 2-6 root cd / && tar czvf /home/backup/stocks.tar.gz -C / var/stock 1>/dev/null
